# scripts/docker/netapi_supervisord.conf
; Supervisor config file for Network API
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;

[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file

[supervisord]
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=5            ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=false               ; start in foreground if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200
environment=C_FORCE_ROOT="true"


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[inet_http_server]
port = 127.0.0.1:9001

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket


;
; Network API daemon settings
;
[program:netapi]
command=/venv/bin/gunicorn --pid /var/run/netapi.pid -c /netapi/gunicorn.conf.py wsgi:application
stdout_logfile=/tmp/networkapi.log
stderr_logfile=/tmp/gunicorn-networkapi_error.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
